# ansible/roles/ros2_node_deploy/templates/ros2-node.service.j2
[Unit]
Description=ROS2 node {{ node_name }} (container)
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=simple
Restart=on-failure
RestartSec=5
TimeoutStopSec=30
# Config dir mounted read-only at node_config_path when set (e.g. /etc/ros2/master2master)
ExecStart=/usr/bin/docker run --rm --name ros2-{{ node_name }} \
  {% if node_config_path is defined and node_config_path is string and node_config_path | length > 0 %}
  -v {{ ros2_node_config_dir }}/{{ node_name }}:{{ node_config_path }}:ro \
  {% endif %}
  {% for env in node_env | default([]) %}
  -e {{ env }} \
  {% endfor %}
  {{ node_extra_args | default('') }} \
  {{ node_image }}
ExecStop=/bin/sh -c '/usr/bin/docker stop -t 25 ros2-{{ node_name }} 2>/dev/null || true'

[Install]
WantedBy=multi-user.target
